# webserver/postgres/jodlgang-db-pwd.service
[Unit]
Before=uwsgi.service
After=postgresql.service

[Service]
ExecStart=/srv/jodlgang/randomify_db_pwd.py

[Install]
WantedBy=default.target